# src/multicache/testbed/stream/server/ffserver.service
[Unit]
Description=ffserver service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/ffserver

[Install]
WantedBy=multi-user.target